# workspace/crontab/root
* * * * * root php /var/www/artisan schedule:run >> /dev/null 2>&1

# Custom cron
* * * * * root echo "Every Minute from laradock" > /var/log/cron.log 2>&1
